# tw-docs-config/systemd/update_tw_docs.service
[Unit]
Description=Update Tumbleweed Community docs

[Service]
ExecStart=/home/lxbuch/bin/update_tw_docs
Type=oneshot
